# infra/systemd/sober-web.service
[Unit]
Description=Sober Web (Frontend + Reverse Proxy)
After=network-online.target sober-api.service
PartOf=sober.target

[Service]
Type=notify
User=sober
Group=sober
ExecStart=/opt/sober/bin/sober-web
EnvironmentFile=/etc/sober/.env
Environment=SOBER_CONFIG=/etc/sober/config.toml
RuntimeDirectory=sober
StateDirectory=sober
Restart=on-failure
RestartSec=5
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes
ReadWritePaths=/opt/sober/data

[Install]
WantedBy=sober.target
