# infra/systemd/sober-agent.service
[Unit]
Description=Sober Agent (gRPC)
After=network-online.target
PartOf=sober.target

[Service]
Type=simple
User=sober
Group=sober
ExecStart=/opt/sober/bin/sober-agent
EnvironmentFile=-/etc/sober/.env
Environment=SOBER_CONFIG=/etc/sober/config.toml
StateDirectory=sober
Restart=on-failure
RestartSec=5
ProtectHome=yes
PrivateTmp=yes
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=sober.target
